# setup-router/ppp-nat/reactive-ppp-interfaces.service
[Unit]
Description=Reactive ppp interfaces
Requires=NetworkManager.service
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=oneshot
Environment=ROUTER_HOME=/home/router
ExecStart=/bin/bash $ROUTER_HOME/ppp-nat/reactive-ppp-interfaces.sh
# RemainAfterExit=yes

[Install]
WantedBy=default.target

# journalctl -e -u reactive-ppp-interfaces.service