# setup-router/sing-box/vbox-switch-haproxy-outbound.service
[Unit]
Description=Vbox auto switch haproxy outbound
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/flock -n /home/router/sing-box/vbox-switch-haproxy-outbound.lock /bin/bash /home/router/sing-box/vbox-switch-haproxy-outbound.sh
# Timer-driven oneshot; flock ensures overlapping triggers discard instantly.

[Install]
WantedBy=default.target

# systemctl enable --user /home/router/sing-box/vbox-switch-haproxy-outbound.service
# systemctl start --user vbox-switch-haproxy-outbound.service
# systemctl status --user vbox-switch-haproxy-outbound.service
