# setup-router/vcs/dockerfile/supervisor-p4d.conf
[program:p4d]
user=root
command=/usr/sbin/p4d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
